# package/lean/luci-app-frpc/root/etc/uci-defaults/luci-frp
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@frpc[-1]
	add ucitrack frpc
	set ucitrack.@frpc[-1].init=frpc
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
